# desktop-user/kwin-monitor-fix/.config/systemd/user/kwin-monitor-fix.service
[Unit]
Description=Second monitor fix

[Service]
ExecStart=%h/.local/bin/kwin-monitor-fix.sh

[Install]
WantedBy=plasma-workspace.target